# 1point3acres_checkin/com.1point3acres.dailyjob.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.1point3acres.dailyjob</string>

    <!-- Run between 12am-1am with a random delay up to 12 hours -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RandomDelay</key>
    <integer>43200</integer> <!-- Max 12 hours -->

    <!-- Retry if missed due to sleep -->
    <key>StartInterval</key>
    <integer>3600</integer> <!-- Every hour -->

    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/tmp/1point3acres_last_run</key>
            <false/>
        </dict>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/wpq/.1point3acres_checkin.sh</string>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/1point3acres_launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/1point3acres_launchd_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>